# ramfs/overlays/BUSYBOX/etc/inittab
::sysinit:/bin/busybox chmod -R +x /etc/init.d/

::sysinit:/etc/init.d/rcS

::sysinit:/sbin/mdev -s
::sysinit:/sbin/mdev -d

ttyS0::respawn:/bin/sh -c 'stty -F /dev/ttyS0 115200; exec /bin/sh </dev/ttyS0 >/dev/ttyS0 2>&1' 
::respawn:/bin/sh

::shutdown:/bin/umount -a -r
::shutdown:/bin/sync